# rc.local
#!/bin/sh -e
#
/etc/init.d/ssh start &
exit 0
